# root/etc/s6-overlay/s6-rc.d/svc-speedtest-tracker/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cd /app/www || exit 1

exec \
    s6-setuidgid abc php artisan queue:work --tries=3 --no-ansi -q
