# root/etc/crontabs/abc
# min   hour    day     month   weekday command
*       *       *       *       *       /usr/bin/php /app/www/artisan schedule:run -q
